# src/remote/libvirtd.socket.in
[Unit]
Description=@name@ local socket
Before=@service@.service
@deps@

[Socket]
ListenStream=@runstatedir@/libvirt/@sockprefix@-sock
Service=@service@.service
SocketMode=0660
SocketUser=root
SocketGroup=libvirt
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
